# .config/systemd/user/camshot.service
[Unit]
Description=Take a screenshot

[Service]
Type=oneshot
ExecStart=/bin/sh -c %h/bin/camshot
